# service_files/nplbam_certbot_renew.timer
[Unit]
Description="Timer that starts the certbot renewal service."

[Timer]
OnBootSec=3min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target